# deploy/leadgen.service
[Unit]
Description=LeadGen AI — FastAPI Backend
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=leadgen
Group=leadgen
WorkingDirectory=/opt/leadgen/backend
EnvironmentFile=/opt/leadgen/.env
ExecStart=/opt/leadgen/.venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --loop uvloop \
    --no-access-log

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=leadgen

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
